# ansible/playbooks/x11/scripts/x11.service
[Unit]
After=network.service

[Service]
KillMode=none
ExecStart=/usr/local/bin/x11-configure.bash

[Install]
WantedBy=default.target
